# setup/systemd/monitor.service
[Unit]
Description=grandma-watcher AI safety assessment loop
After=network-online.target go2rtc.service
Wants=network-online.target
Requires=go2rtc.service

[Service]
Type=simple
User=pi
Group=pi

ExecStart=/home/pi/eldercare/.venv/bin/python /home/pi/eldercare/monitor.py

WorkingDirectory=/home/pi/eldercare

# Restart policy
Restart=on-failure
RestartSec=5

# Resource limits (Pi 5 8GB)
MemoryMax=512M

# Environment
Environment=HOME=/home/pi

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grandma-monitor

[Install]
WantedBy=multi-user.target
